# agora-ui.service
[Unit]
Description=Agora UI - Browser chat interface for Agora messaging
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stefan
WorkingDirectory=/home/stefan/agora-ui
ExecStart=/home/stefan/.nvm/versions/node/v24.14.0/bin/node dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
